# docker/supervisord.conf
# Supervisor 配置文件
# 用于管理 Nginx 和 FastAPI 后端进程

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

# Nginx 进程配置
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log
priority=1
user=appuser
environment=PYTHONUNBUFFERED="1"

# FastAPI 后端进程配置
[program:backend]
command=/home/appuser/.local/bin/uvicorn backend.main_ai_edu:app --host 0.0.0.0 --port 8000 --workers 4
directory=/app/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/backend/error.log
stdout_logfile=/var/log/backend/access.log
priority=2
user=appuser
environment=PYTHONPATH="/app/backend",APP_ENV="production",LOG_LEVEL="info"
stopasgroup=true
killasgroup=true
stopsignal=TERM
numprocs=1
process_name=%(program_name)s_%(process_num)02d
